# extern/etc/cron.d/cloudflare-ip-whitelist-sync
SHELL=/bin/bash
PATH=/usr/sbin:/usr/bin:/sbin:/bin

# Daily Cloudflare whitelist sync
30 4 * * * root /bin/bash /etc/cloudflare/ip_whitelist_sync.sh >> /var/log/cloudflare-ip-whitelist-sync.log 2>&1
