# extern/etc/systemd/system/fanic.service
[Unit]
Description=FANIC web service
After=network.target

[Service]
Type=simple
User=fanic
Group=www-data
SupplementaryGroups=fanic
WorkingDirectory=/opt/fanic
Environment=PATH=/opt/fanic/.venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/opt/fanic/src
EnvironmentFile=-/etc/fanic/fanic.env
Environment=FANIC_DATA_DIR=/mnt/storage
Environment=FANIC_DB_PATH=/opt/fanic/runtime/fanic.db
Environment=FANIC_LOG_PATH_TEMPLATE=/var/log/fanic/%%TIMESTAMP%%.log
Environment=PYTHONDONTWRITEBYTECODE=1
# Optional secure secret file pointers (recommended for production):
# Environment=FANIC_SESSION_SECRET_FILE=/etc/fanic/credentials/session_secret
# Environment=FANIC_AUTH0_CLIENT_SECRET_FILE=/etc/fanic/credentials/auth0_client_secret
# Environment=FANIC_AUTH0_AUDIENCE_FILE=/etc/fanic/credentials/auth0_audience
# Environment=FANIC_ALPHA_INVITE_CODES_CSV_FILE=/etc/fanic/credentials/alpha_invite_codes_csv
UMask=0027
PrivateTmp=true
ProtectSystem=strict
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictRealtime=true
RestrictNamespaces=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible
ProcSubset=pid
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
TasksMax=512
MemoryMax=8G
RuntimeDirectory=fanic
RuntimeDirectoryMode=0750
ReadWritePaths=/mnt/storage /opt/fanic/runtime /var/log/fanic /run/fanic
ExecStart=/usr/bin/aa-exec -p fanic -- /opt/fanic/.venv/bin/python -m fanic.main serve --unix-socket /run/fanic/fanic.sock --unix-socket-perms 660
Restart=on-failure
RestartSec=5
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
